# systemd/user/repo-checkin.timer
[Unit]
Description=Daily Git Repo Checkin

[Timer]
OnCalendar=*-*-* 03:30
Persistent=true

[Install]
WantedBy=timers.target
